# files/system/etc/profile.d/pixi-path.sh
# pixi-path.sh
# Always add pixi bin to PATH so it works immediately after first install
export PATH="$HOME/.pixi/bin:$PATH"